# deployment/supervisor/microblog.conf
[program:blog]
command=/home/ubuntu/blog/venv/bin/gunicorn -b localhost:8000 -w 4 blog:app
directory=/home/ubuntu/blog
user=ubuntu
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
