# deployment/supervisor/microblog-tasks.conf
[program:blog-tasks]
command=/home/ubuntu/blog/venv/bin/rq worker blog-tasks
numprocs=1
directory=/home/ubuntu/blog
user=ubuntu
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
